# templates/podman.service.j2
[Unit]
Description={{ description }}

[Service]
Restart=always
ExecStart=/usr/bin/podman start -a {{ name }}
ExecStop=/usr/bin/podman stop -t 2 {{ name }}

[Install]
WantedBy=local.target